# host/launchd/com.majordomus.taskrouter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.majordomus.taskrouter — Majordomus headless Task Router host (always-on LaunchAgent)

  TEMPLATE — secrets are NEVER committed here.
  Run  bash host/provision.sh --inject-secrets  to patch this file in-place
  with real values from .env before loading it with launchctl.

  Load (once per session or boot):
    launchctl bootstrap gui/$(id -u) \
      ~/GitHub/Majordomos/host/launchd/com.majordomus.taskrouter.plist

  Unload:
    launchctl bootout gui/$(id -u)/com.majordomus.taskrouter

  Logs:  ~/Library/Logs/majordomus/taskrouter.log / taskrouter.err
-->
<plist version="1.0">
<dict>

  <!-- Identity -->
  <key>Label</key>
  <string>com.majordomus.taskrouter</string>

  <!-- Run start-majordomos.sh through a bash login shell so nvm / Homebrew
       PATH is inherited.  The script also sources .env; EnvironmentVariables
       below provides an additional explicit layer for the key secrets. -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-l</string>
    <!-- patched by provision.sh --inject-secrets -->
    <string>__INJECT_AT_PROVISION__/start-majordomos.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <!-- patched by provision.sh --inject-secrets -->
  <string>__INJECT_AT_PROVISION__</string>

  <!-- Secrets injected at provision time — never committed with real values.
       All values below are placeholder sentinels; provision.sh replaces them. -->
  <key>EnvironmentVariables</key>
  <dict>
    <!-- Core auth (G1 — headless, no interactive login) -->
    <key>ANTHROPIC_API_KEY</key>
    <string>__INJECT_AT_PROVISION__</string>

    <!-- Task Router API key (gates the REST surface) -->
    <key>TASK_ROUTER_API_KEY</key>
    <string>__INJECT_AT_PROVISION__</string>

    <!-- Home Assistant MCP SSE endpoint + bearer token (used by .mcp.json
         and by the v2 ha-bridge loopback executor) -->
    <key>HA_BASE_URL</key>
    <string>__INJECT_AT_PROVISION__</string>
    <key>HA_TOKEN</key>
    <string>__INJECT_AT_PROVISION__</string>

    <!-- Expands ${HOME} for node-pty / nvm lookups inside the supervisor -->
    <key>HOME</key>
    <string>__INJECT_AT_PROVISION__</string>
  </dict>

  <!-- Restart on crash; ThrottleInterval prevents hot-loop restart storms
       (10 s minimum between restarts — matches watchdog.ts). -->
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>RunAtLoad</key>
  <true/>

  <!-- Logs — patched by provision.sh --inject-secrets to an absolute path -->
  <key>StandardOutPath</key>
  <string>__INJECT_AT_PROVISION__/Library/Logs/majordomus/taskrouter.log</string>
  <key>StandardErrorPath</key>
  <string>__INJECT_AT_PROVISION__/Library/Logs/majordomus/taskrouter.err</string>

</dict>
</plist>
